# MultimodelRagGroq/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid

; FastAPI — public on port 7860 (HF Spaces default)
[program:api]
command=uvicorn app.main:app --host 0.0.0.0 --port 7860 --workers 1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/api.log
stderr_logfile=/var/log/supervisor/api_err.log
environment=UPLOAD_DIR="/data/uploads"

; Celery worker — concurrency 1 to stay within free-tier RAM
[program:worker]
command=celery -A app.workers.celery_app worker --loglevel=info --concurrency=1
directory=/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/worker.log
stderr_logfile=/var/log/supervisor/worker_err.log
environment=UPLOAD_DIR="/data/uploads"
